# systemd/user/weston-kiosk.socket
[Unit]
Description=Weston kiosk socket

[Socket]
ListenStream=%t/wayland-0
